# objects/systemd-units/units/tmux-refresh-clients.service
[Unit]
Description=Tmux clients auto-refresh daemon
Wants=init-tmux-refresh-clients.service
After=init-tmux-refresh-clients.service

[Service]
ExecStart=/usr/bin/bash %h/scripts/services/tmux-refresh-clients.sh
Restart=on-failure

[Install]
WantedBy=default.target

