# objects/systemd-units/units/update-calendar.service
[Unit]
Description=Calendar update job

[Service]
Type=oneshot
ExecStart=/usr/bin/touch %h/.when/calendar

