# special/systemd/update-calendar.timer
[Unit]
Description=Calendar update job timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

