# special/systemd/update-calendar.service
[Unit]
Description=Calendar update job

[Service]
Type=oneshot
ExecStart=/usr/bin/touch %h/Org/calendar

